# addon/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup when nginx stops
# ==============================================================================

bashio::log.info "Nginx service stopped"
