# addon/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start nginx
# ==============================================================================

bashio::log.info "Starting nginx..."

# Wait for backend to be ready
sleep 2

exec nginx -g "daemon off;"
