# addon/rootfs/etc/services.d/backend/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup when backend stops
# ==============================================================================

bashio::log.info "Backend service stopped"
